# prism_tier2/prism_tier2/systemd/prism-watcher.service
[Unit]
Description=PRISM File System Watcher
After=network.target redis.service prism-worker-classify.service
Requires=redis.service

[Service]
Type=simple
User=prism
Group=prism
WorkingDirectory=/opt/prism
EnvironmentFile=/opt/prism/.env
Environment=PYTHONPATH=/opt/prism
ExecStart=/opt/prism/venv/bin/python watcher.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/opt/prism/logs/watcher.log
StandardError=append:/opt/prism/logs/watcher.log

[Install]
WantedBy=multi-user.target
